# scripts/every20min_generateCrousMenus.cron
*/20 * * * * root /path/to/generateCrousMenus.sh
